# ansible/playbooks/templates/filesystem-livesync.service.j2
[Unit]
Description=Filesystem LiveSync for Obsidian
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User={{ ansible_user }}
WorkingDirectory=/home/{{ ansible_user }}/filesystem-livesync
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
